# infra/fetchium-healthcheck.timer
[Unit]
Description=Run Fetchium API health watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
RandomizedDelaySec=15s
Persistent=true
Unit=fetchium-healthcheck.service

[Install]
WantedBy=timers.target
