# infra/fetchium-secretops-refresh.service
[Unit]
Description=Fetchium SecretOps refresh fan-out
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart fetchium-api fetchium-admin fetchium-dashboard fetchium-web

[Install]
WantedBy=multi-user.target
