# infra/fetchium-healthcheck.service
[Unit]
Description=Fetchium API health watchdog
After=network.target fetchium-api.service

[Service]
Type=oneshot
EnvironmentFile=-***REMOVED***
ExecStart=/home/echo/projects/Fetchium/scripts/fetchium-healthcheck.sh
